# c7-foreman/system/foreman-install.target
[Unit]
Description=Install Foreman And Exit
Requires=multi-user.target
